# service/lliurex-unattended-upgrades.service
[Unit]
Description="Install packages from custom component for security"

[Service]
ExecStart=/usr/bin/lliurex-unattended-upgrades
